# deploy/systemd/contradiction-radar-watchdog.timer
[Unit]
Description=Run the Contradiction Radar health watchdog every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=20s
Persistent=true

[Install]
WantedBy=timers.target
